# etc/cron.d/fetchmail
0/5 * * * * fetchmail -v -f /etc/request-tracker4/fetchmail
